# debian/systemd/runit.service
[Unit]
Description=Runit service supervision

[Service]
ExecStart=/usr/sbin/runsvdir-start
Restart=always

[Install]
WantedBy=multi-user.target
